# root/etc/uci-defaults/90_dae
#!/bin/sh
# UCI defaults for luci-app-dae — runs on first boot

# Ensure the config section exists
if ! uci get dae.settings >/dev/null 2>&1; then
    uci set dae.settings=settings
    uci set dae.settings.enabled='1'
    uci set dae.settings.config_file='/etc/dae/config.dae'
    uci set dae.settings.disable_waiting_network='0'
    uci set dae.settings.enable_dae_log='1'
    uci set dae.settings.logfile_maxsize='1'
    uci commit dae
fi

mkdir -p /var/log/dae
touch /var/log/dae/dae.log
chmod 644 /var/log/dae/dae.log

mkdir -p /etc/dae
touch /etc/dae/config.dae
chmod 600 /etc/dae/config.dae

exit 0
